# config_files/dump_grabber
#!/sbin/runscript

depend() {
  need net
  use dns localmount chaosc
  after bootmisc
  provide dump_grabber
}


start() {
  ebegin "starting dump_grabber"
  start-stop-daemon --start --pidfile /var/run/dump_grabber.pid --make-pidfile --user stefan --group stefan --background --exec env DISPLAY=:0 /usr/bin/dump_grabber
  eend $?
}

stop() {
  ebegin "stopping dump_grabber"
  start-stop-daemon --stop --quiet --pidfile /var/run/dump_grabber.pid
  eend $?
}
