# config_files/ekgplotter
#!/sbin/runscript

depend() {
  need net
  use dns localmount chaosc
  after bootmisc
  provide ekgplotter
}


start() {
  ebegin "starting ekgplotter"
  start-stop-daemon --start --pidfile /var/run/ekgplotter.pid --make-pidfile --user stefan --group stefan --background --exec env DISPLAY=:0 /usr/bin/ekgplotter
  eend $?
}

stop() {
  ebegin "stopping ekgplotter"
  start-stop-daemon --stop --quiet --pidfile /var/run/ekgplotter.pid
  eend $?
}
